# salt/states/codered/autoupdate/files/codered-update.timer
[Unit]
Description=CodeRed AI Sensor Auto-Update Timer

[Timer]
OnBootSec=15min
OnUnitActiveSec=6h
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
